# templates/etc/init/kibana-RedHat.conf.erb
description "start and stop kibana"
author "Gareth Rushgrove"

start on runlevel [23]
stop on runlevel [06]

respawn
respawn limit 5 30

chdir /opt/kibana

console output
script
exec su -s /bin/sh -c 'bundle exec ruby kibana.rb' kibana
end script
