# templates/etc/init/kibana-RedHat.conf.erb
description "start and stop kibana"
author "Gareth Rushgrove"

start on runlevel [23]
stop on runlevel [06]

respawn
respawn limit 5 30

chdir /opt/kibana
setuid kibana
setgid kibana

console output
script
exec bundle exec ruby kibana.rb
end script
